# meta-customnux-dashboard/recipes-support/ota-status-check/files/ota-status-check.path
[Unit]
Description=file creation observer: /tmp/swu.status

[Path]
PathExists=/tmp/swu.status
Unit=ota-status-check.service

[Install]
WantedBy=multi-user.target
