# meta-customnux-dashboard/recipes-support/ota-status-check/files/ota-status-check.service
[Unit]
Description=SWUpdate post handler

[Service]
Type=oneshot
ExecStart=/usr/bin/ota-status-check.sh
